# setup/client/etc/rc.local
xinput --set-prop 'FT5406 memory based driver' 'Coordinate Transformation Matrix'  0 1 0 -1 0 1 0 0 1
echo 72 > /sys/class/backlight/rpi_backlight/brightness
